# slider.initd
#!/bin/sh /etc/rc.common

START=97
SERVICE_NAME=slider
SERVICE_DAEMONIZE=1
SERVICE_WRITE_PID=1

start() {
  service_check /bin/sh && return 1
  service_start /etc/slider daemon
}

stop() {
  service_stop /bin/sh
}

